# deploy/ai.openclaw.dashboard.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.dashboard</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/tianshuai/.openclaw/workspace/dashboard-live/start.sh</string>
      <string>8888</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>OPENCLAW_BASE_DIR</key>
      <string>/Users/tianshuai/.openclaw</string>
      <key>OPENCLAW_WORKSPACE</key>
      <string>/Users/tianshuai/.openclaw/workspace</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/tianshuai/.openclaw/workspace/dashboard-live</string>
    <key>StandardOutPath</key>
    <string>/Users/tianshuai/.openclaw/workspace/dashboard-live/.tmp/dashboard.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tianshuai/.openclaw/workspace/dashboard-live/.tmp/dashboard.launchd.err.log</string>
  </dict>
</plist>
